# SYSTEM/etc/systemd/system/kanetKB.service
[Unit]
Description=kanetKB service.

[Service]
Type=simple
ExecStart=/bin/bash /opt/VitaniveKB/bin/ServiceD.sh

[Install]
WantedBy=multi-user.target
